# templates/graphital.conf.erb
# This is a upstart configfile
#
# <%= @graphite_host %>
 
description     "send system vitals to graphite"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /usr/bin/ruby /opt/graphital/graphital.rb
